# PKGBUILDs/mariadb-meta/mariadb-meta.install
post_install() {
   mariadb-install-db --user=mysql --basedir=/usr --datadir=/var/lib/mysql
   systemctl enable --now mariadb.service
}
pre_remove() {
}
