# PKGBUILDs/gnome-session-minimal/gnome-session-minimal.install
post_install() {
   systemctl enable NetworkManager.service
   systemctl enable bluetooth.service
}
