# PKGBUILDs/apache-meta/apache-meta.install
post_install() {
   systemctl enable httpd.service
}
