# PKGBUILDs/gdm-meta/gdm-meta.install
post_install() {
   systemctl enable gdm.service
}
